# systemd/tech-demo-client@.service.d/dependencies.conf
[Unit]
After=tech-demo-config-helper.socket
Wants=tech-demo-config-helper.socket